# files/deepin-notification.service.in
[Unit]
Description=DDE notification service

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart=$$PREFIX/lib/deepin-notifications
